# Crontab/crontab
1 * * * * echo "Hello world!" >> /var/log/syslog
